# systemd/.config/systemd/user/git-maintenance@.timer
[Unit]
Description=Prefetch %i.git

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySeconds=600

[Install]
WantedBy=timers.target

